# mono-devel-list/attachments/20040225/854fa376/attachment.obj
#! /bin/bash
#
# chkconfig: 2345 99 01
# description: Provides ASP.NET functionality for Apache
#

. /etc/rc.d/init.d/functions

case "$1" in
  start)
	echo -n "Starting ASP.NET..."
	su - aspnet /home/aspnet/bin/start_aspnet
	success "aspnet startup"
	echo
	;;

  stop)
	echo -n "Stopping ASP.NET..."
	su - aspnet /home/aspnet/bin/stop_aspnet
	success "aspnet cleanup"
	echo
	;;

  restart)
	$0 stop
	$0 start
	;;

  *)
	echo $"Usage: $0 {start|stop|restart}"
	exit 1
esac

exit 0

